# Actividades SystemD/enviaCorreo.service
[Unit]
Description=Servicio Fallido

[Service]
ExecStart=$HOME/enviaCorreo.sh

[Install]
WantedBy=multi-user.target
